# systemd/backlight.service
[Unit]
Description=Backlight

[Service]
Type=simple
ExecStart=/etc/systemd/system/backlight

[Install]
WantedBy=multi-user.target
